# deploy/systemd/agents-chat-backup.service
[Unit]
Description=Agents Chat local backup job
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/ops/run-backups.sh
